# configs/uwsgi.service
[Unit]
Description=uWSGI Emperor service

[Service]
ExecStart=/home/lamp/.uwsgi/python35/bin/uwsgi --emperor /etc/uwsgi/sites
Restart=always
KillSignal=SIGQUIT
Type=notify
NotifyAccess=all

[Install]
WantedBy=multi-user.target

